# supervisord.conf
[supervisord]
nodaemon=true

[program:chrome]
command=/dockerstartup/vnc_startup.sh --skip
priority=10
startsecs=5
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
